# supervisord.conf
[supervisord]
#nodaemon=true
directory=/muru-proxy

[program:webp-server]
command=./webp-server